# mgmt/ceph-deploy.conf
[program:ceph-deploy]
command=bash -c 'sleep 20 && sudo -u ceph -i /home/ceph/bootceph.sh'
autorestart=false
redirect_stderr=true
