# mgmt/ceph-deploy.conf
[program:ceph-deploy]
command=bash -c 'sleep 5 && sudo -u ceph -i /home/ceph/bootstrap.sh'
autorestart=false
redirect_stderr=true